# infra/plex-hwaccel/plex-hwaccel.service
[Unit]
Description=Plex Hardware Acceleration Device Access
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/plex-hwaccel.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
